# glinet_gl-mt300n-v2/etc/uci-defaults/982-dhcp
#!/bin/sh

set -e

uci set dhcp.@dnsmasq[0].local="/internal/"
uci set dhcp.@dnsmasq[0].domain="internal"
uci commit dhcp
