# glinet_gl-ar750/etc/uci-defaults/983-ssh
#!/bin/sh

set -e

uci set dropbear.@dropbear[0].RootPasswordAuth="off"
uci set dropbear.@dropbear[0].Interface="lan"
uci commit dropbear
