# glinet_gl-mt1300/etc/uci-defaults/985-network
#!/bin/sh

set -e

uci set network.lan.ipaddr="192.168.8.1"
uci commit network
